# debian/lokkit-chain.postinst
#!/bin/bash

#update-rc.d administra defaults
